# packaging/debian/prerm
#!/bin/sh
set -e

if command -v systemctl >/dev/null 2>&1; then
  systemctl disable --now infinifix-thresholds.service >/dev/null 2>&1 || true
fi

if command -v update-desktop-database >/dev/null 2>&1; then
  update-desktop-database /usr/share/applications >/dev/null 2>&1 || true
fi

exit 0
